# ipod-listener.service
[Unit]
Description=iPod USB listener
After=network.target

[Service]
Type=simple
ExecStart=/opt/ipod-dock/.venv/bin/python -m ipod_sync.udev_listener
WorkingDirectory=/opt/ipod-dock
Restart=on-failure
User=ipod
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
