# ipod-watcher.service
[Unit]
Description=iPod sync queue watcher
After=network.target

[Service]
Type=simple
ExecStart=/opt/ipod-dock/.venv/bin/python -m ipod_sync.watcher
WorkingDirectory=/opt/ipod-dock
Restart=on-failure
User=ipod
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
