# ipod-api.service
[Unit]
Description=iPod FastAPI server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/usr/bin/python3 -m ipod_sync.app
Restart=on-failure
User=ipod
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
