# ipod-watcher.service
[Unit]
Description=iPod sync queue watcher
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m ipod_sync.watcher
WorkingDirectory=%h
Restart=on-failure
User=ipod
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
